# supporting-materials/systemd-units-overrides/etc/systemd/system/haveged.service.d/01-honour-default.conf
# systemd crap is hardcoded and does not respect envs
[Service]
EnvironmentFile=/etc/default/haveged
ExecStart=
ExecStart=/usr/sbin/haveged --Foreground --verbose=1 ${DAEMON_ARGS}
